# scripts/security/systemd/credential-rotation-monthly.timer
[Unit]
Description=Monthly Credential Rotation Timer (Tier 2)
Documentation=file:///home/gerald/IdeaProjects/Datamancy/scripts/security/README.md
Requires=credential-rotation-monthly.service

[Timer]
# Run on the 1st of every month at 4:00 AM
OnCalendar=*-*-01 04:00:00

# If system was off, run on next boot
Persistent=true

# Add random delay of 0-5 minutes
RandomizedDelaySec=300

AccuracySec=1min

[Install]
WantedBy=timers.target
