# scripts/security/systemd/credential-rotation.service
[Unit]
Description=Weekly Credential Rotation for Datamancy Stack
After=network.target docker.service
Requires=docker.service
Documentation=file:///home/gerald/IdeaProjects/Datamancy/scripts/security/README.md

[Service]
Type=oneshot
User=gerald
Group=gerald
WorkingDirectory=/home/gerald/IdeaProjects/Datamancy

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/snap/bin"
Environment="KOTLIN_HOME=/snap/kotlin/current"

# Pre-flight checks
ExecStartPre=/usr/bin/docker ps --quiet --filter name=postgres
ExecStartPre=/usr/bin/bash -c 'test -f /home/gerald/IdeaProjects/Datamancy/.env'

# Main rotation script
ExecStart=/usr/bin/kotlin /home/gerald/IdeaProjects/Datamancy/scripts/security/weekly-rotation.main.kts --execute

# Post-rotation health check
ExecStartPost=/usr/bin/kotlin /home/gerald/IdeaProjects/Datamancy/scripts/security/lib/health-check.main.kts --execute

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=credential-rotation

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/gerald/IdeaProjects/Datamancy/secrets

# Timeouts
TimeoutStartSec=900
TimeoutStopSec=60

# Restart policy (don't restart on failure, let timer handle it)
Restart=no

[Install]
WantedBy=multi-user.target
