# scripts/security/systemd/credential-rotation-monthly.service
[Unit]
Description=Monthly Credential Rotation for Datamancy Stack (Tier 2)
After=network.target docker.service
Requires=docker.service
Documentation=file:///home/gerald/IdeaProjects/Datamancy/scripts/security/README.md

[Service]
Type=oneshot
User=gerald
Group=gerald
WorkingDirectory=/home/gerald/IdeaProjects/Datamancy

Environment="PATH=/usr/local/bin:/usr/bin:/bin:/snap/bin"
Environment="KOTLIN_HOME=/snap/kotlin/current"

ExecStartPre=/usr/bin/docker ps --quiet --filter name=postgres
ExecStartPre=/usr/bin/bash -c 'test -f /home/gerald/IdeaProjects/Datamancy/.env'

ExecStart=/usr/bin/kotlin /home/gerald/IdeaProjects/Datamancy/scripts/security/monthly-rotation.main.kts --execute

StandardOutput=journal
StandardError=journal
SyslogIdentifier=credential-rotation-monthly

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/gerald/IdeaProjects/Datamancy/secrets

TimeoutStartSec=900
TimeoutStopSec=60

Restart=no

[Install]
WantedBy=multi-user.target
